# bot.conf
[program:DemaTgBot]
directory=/home/ubuntu/DemaTgBot
command=python3 bot.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu",USER="ubuntu"
stderr_logfile=/home/ubuntu/DemaTgBot/logfile_err.log
stdout_logfile=/home/ubuntu/DemaTgBot/logfile.log